# com.mcp-workshop.servers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this service -->
    <key>Label</key>
    <string>com.mcp-workshop.servers</string>
    
    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/spencer/repos/mcp-server/scripts/launchd-start.sh</string>
    </array>
    
    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep alive - restart if it exits -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <!-- Restart throttling - wait 10 seconds between restarts -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/spencer/repos/mcp-server</string>
    
    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/spencer/Library/Logs/MCP/mcp-servers.out.log</string>
    
    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/spencer/Library/Logs/MCP/mcp-servers.err.log</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
    
    <!-- Process type -->
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <!-- Nice value (priority) -->
    <key>Nice</key>
    <integer>0</integer>
    
    <!-- Soft resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
</dict>
</plist>